# root/etc/s6-overlay/s6-rc.d/init-envfile/run
#!/command/with-contenv bash

mapfile -t file_vars < <(compgen -e | grep '^FILE__')

if [ ${#file_vars[@]} -eq 0 ]; then
  echo "[init-envfile] No FILE__ variables found, skipping..."
  exit 0
fi

for var_name in "${file_vars[@]}"; do
  var_value="${!var_name}"
  target_var="${var_name#FILE__}"

  if [ -f "$var_value" ]; then
    printf '%s' "$(cat "$var_value")" > /run/s6/container_environment/"$target_var"
    echo "[init-envfile] Set $target_var from $var_value"
    rm -f /run/s6/container_environment/"$var_name"
    echo "[init-envfile] Unset $var_name"
  else
    echo "[init-envfile] Warning: $var_name points to $var_value which does not exist, skipping"
  fi
done
