# root/etc/s6-overlay/s6-rc.d/svc-cron-log/run
#!/command/with-contenv bash

while IFS= read -r line; do
  printf '[%s] %s\n' "$(date '+%Y-%m-%dT%H:%M:%S')" "$line"
done
