# root/etc/s6-overlay/s6-rc.d/svc-cron/finish
#!/bin/bash

if test "$1" -eq 256 ; then
  e=$((128 + $2))
else
  e="$1"
fi

echo "Received exit code $e."
echo "$e" > /run/s6-linux-init-container-results/exitcode